# dnsdist/rootfs/etc/services.d/recursor/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid abc \
    /usr/bin/dnsdist \
    --supervised \
    --disable-syslog \
    --config /etc/dnsdist/dnsdist.yml
