# lightningstream/rootfs/etc/services.d/lightningstream/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid abc \
    /usr/local/bin/lightningstream \
    --config="/lightningstream.yml" \
    sync
